# salt/borda/borda.conf
# broker - upstart job file

description "borda measurement collector"
author "Ox Cart <ox@getlantern.org>"

# Stanzas
#
# Stanzas control when and how a process is started and stopped
# See a list of stanzas here: http://upstart.ubuntu.com/wiki/Stanzas#respawn

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [016]

# Automatically restart process if crashed
respawn

# Run as lantern
setuid lantern

# Run in lantern's home directory
chdir /home/lantern

# We've observed we need to raise the limit on open file descriptors
limit nofile 1024768 1024768

# Start the process, piping output prepended with a prefix to syslog
# Note - the privkey.pem and fullchain.pem have to be deployed manually at the
#        moment for security reasons.
# Note - the file authkey has to be manually populated with the authentication
#        key.
exec /home/lantern/borda -config /home/lantern/config.ini -pkfile privkey.pem -certfile fullchain.pem -stderrthreshold WARNING 2>&1 | logger -t borda
