# salt/cloudmaster/destroy.conf
# destroy - upstart job file

description "destroy"
author "<aranhoide@getlantern.org>"

# Stanzas
#
# Stanzas control when and how a process is started and stopped
# See a list of stanzas here: http://upstart.ubuntu.com/wiki/Stanzas#respawn

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [016]

# Automatically restart process if crashed
respawn

# Terminate with SIGINT instead of SIGKILL
kill signal SIGINT

# Wait 15 seconds before sending SIGKILL if SIGINT doesn't work
kill timeout 15

# Run as root
setuid root

# Run in lantern's home directory
chdir /home/lantern

# We've observed we need to raise the limit on open file descriptors
limit nofile 1024768 1024768

# Set environment variables for configuration and authentication
env CM="{{ pillar['cloudmaster_name'] }}"
export CM
env DO_TOKEN="{{ pillar['do_token'] }}"
export DO_TOKEN
env VULTR_APIKEY="{{ pillar['vultr_apikey'] }}"
export VULTR_APIKEY
env REDIS_URL="{{ pillar['cfgsrv_redis_url'] }}"
export REDIS_URL
env PYTHONPATH="/usr/local/lib/pylib"
export PYTHONPATH

# Start the process, piping output prepended with a prefix to syslog
exec /usr/bin/destroy.py 2>&1 | logger -t destroy
