# salt/config_server/config-server.conf
description "config-server http daemon"
author "Alan Hoyte (aranhoide@gmail.com)"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

chdir /home/lantern
limit nofile 128000 128000
exec su lantern -c "java -Xms800m -Xmx1600m -jar /home/lantern/config-server.jar 2>&1 | logger -t config-server" 2>&1 | logger -t config-server
